# etc/systemd/system/rabbitmq_exporter.service
[Unit]
Description=Rabbitmq Exporter
After=docker.service
Requires=docker.service

[Service]
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker stop rabbitmq_exporter
ExecStartPre=-/usr/bin/docker rm rabbitmq_exporter
ExecStartPre=/usr/bin/docker pull kbudde/rabbitmq-exporter
ExecStart=docker run --rm --name rabbitmq-exporter -e RABBIT_URL=http://localhost:15672 --network=host kbudde/rabbitmq-exporter
Restart=always

[Install]
WantedBy=multi-user.target
